# systemd-linux/quip-network-node.service
[Unit]
Description=Quip Network Node
After=network.target
Wants=network.target

[Service]
Type=simple
User=quip
Group=quip
Environment=PATH=/opt/quip/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/quip/lib/python3.13/site-packages
ExecStart=/opt/quip/bin/quip-network-node --config /etc/quip.network/config.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quip-network-node

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/quip.network /var/lib/quip.network
ProtectKernelTunables=yes
ProtectControlGroups=yes

# Resource limits
#MemoryMax=32G
#CPUQuota=800%

[Install]
WantedBy=multi-user.target
